# connection_troubleshooting/whatsapp-keepalive.service
[Unit]
Description=WhatsApp Gateway Keep-Alive Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw/.openclaw/workspace
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /opt/openclaw/keepalive/keep_alive.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /var/tmp /home/openclaw/.openclaw

[Install]
WantedBy=multi-user.target